# code/collect/collect_twitter/v1-streaming/supervisor/midterm_stream_candidates.conf
[program:midterm_stream_candidates]]
command = **PATH TO PYTHON3** /**REPO DIR**/collect/collect_twitter/v1-streaming/candidates-twitter-streamer-V1.py -c /**REPO DIR**/config.ini     ; the program (relative uses PATH, can take args)
directory = /home/**YOUR USERNAME**/midterm2022/                                    ; directory to cwd to before exec (def no cwd)
autostart = true                                                                    ; start at supervisord start (default: true)
autorestart = true                                                                  ; when to restart if exited after running (def: unexpected)
user = **YOUR USERNAME**                                                            ; setuid to this UNIX account to run the program
stderr_logfile = /var/log/supervisor/midterm_stream_candidates_err.log               ; stderr log path, NONE for none; default AUTO
stdout_logfile = /var/log/supervisor/midterm_stream_candidates.log                   ; stdout log path, NONE for none; default AUTO